# niri/.config/systemd/user/auto-theme.service
[Unit]
Description=Auto theme switcher based on sunrise/sunset
After=network.target niri.service

[Service]
Type=oneshot
ExecStart=/home/ashzak/.config/niri/auto-theme.py
Environment=DISPLAY=:0
