# systemd/.config/systemd/user/auto-theme.service
[Unit]
Description=Auto theme switcher based on sunrise/sunset
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.config/auto-theme.py
Environment=DISPLAY=:0
