# labwc/.config/systemd/user/auto-theme.service
[Unit]
Description=Auto theme switcher based on sunrise/sunset
After=network.target labwc.service

[Service]
Type=oneshot
ExecStart=/home/ashzak/.config/labwc/auto-theme.sh
Environment=DISPLAY=:0
